# templates/minecraft.socket.j2
[Socket]
ListenFIFO={{ minecraft_console_fifo }}
SocketMode={{ minecraft_console_fifo_mode }}
